# roles/gui/files/systemd/user/pulseaudio-dlna.service
[Unit]
Description=Pulseaudio dlna

[Service]
Type=simple
ExecStart=/usr/bin/pulseaudio-dlna

[Install]
WantedBy=default.target
